# deploy/systemd/aa-scrape.service
[Unit]
Description=AA meetings daily scrape and backup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=flaskuser
Group=flaskuser
WorkingDirectory=/var/www/aa
EnvironmentFile=/var/www/aa/deploy/aa.env
ExecStart=/var/www/aa/.venv/bin/python /var/www/aa/scripts/scrape_and_backup.py
StandardOutput=journal
StandardError=journal
