# deploy/systemd/aa-web.service
[Unit]
Description=AA meetings web service
After=network.target

[Service]
Type=simple
User=flaskuser
Group=flaskuser
WorkingDirectory=/var/www/aa
EnvironmentFile=/var/www/aa/deploy/aa.env
ExecStart=/var/www/aa/.venv/bin/gunicorn --workers 3 --bind 127.0.0.1:8008 wsgi:app
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
